# scripts/pakemon.service
[Unit]
Description=Pakémon

[Service]
ExecStart=cd /opt/pakemon/ && node index.js
Restart=always
RestartSec=5s

[Install]
WantedBy=basic.target
